# templates/dl-unsplash-wallpaper.service.j2

[Unit]
Description=Download unsplash wallpaper

[Service]
ExecStart={{ dl_wallpaper_script_path }}
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
