# templates/rclone_sync.service.j2
[Unit]
Description=Rclone sync service of {{ item.sync_path }} to {{ item.sync_remote }}

[Service]
ExecStart={{ home }}/.local/bin/rclone_sync_{{ item.name }}.sh

[Install]
WantedBy=default.target
