# templates/umount_veracrypt_unit.service.j2
[Unit]
Description=Unmount all veracrypt devices on shutdown

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=unmount_veracrypt_drives

[Install]
WantedBy=multi-user.target
